# hostberry-web.service
[Unit]
Description=HostBerry Web Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/hostberry
ExecStart=/opt/hostberry/venv/bin/gunicorn --workers 3 --bind 0.0.0.0:80 --access-logfile /var/log/hostberry/access.log --error-logfile /var/log/hostberry/error.log app:app
Restart=always
RestartSec=10
Environment="FLASK_APP=app.py"
Environment="FLASK_ENV=production"
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target